# docker/supervisor/laravel/example.conf
[group:main]
programs=main

[program:main]
command=php /var/www/artisan queue:work
numprocs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
